# IBKR_trading/systemd/ibkr-trading.service
[Unit]
Description=IBKR Live Trading System
Documentation=https://github.com/your-org/trader
After=network.target
Wants=network.target

[Service]
Type=simple
User=trader
Group=trader
WorkingDirectory=/home/Jennifer/secure/trader/IBKR_trading
ExecStart=/home/Jennifer/miniconda3/bin/python live_trading/ibkr_live_exec.py --config config/ibkr_live.yaml --symbols AAPL,MSFT,TSLA,GOOGL,AMZN
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Environment
Environment=PYTHONPATH=/home/Jennifer/secure/trader
Environment=PYTHONUNBUFFERED=1
Environment=IBKR_TRADING_MODE=live

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/Jennifer/secure/trader/IBKR_trading
ReadWritePaths=/home/Jennifer/secure/trader/logs
ReadWritePaths=/home/Jennifer/secure/trader/state

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=8G
CPUQuota=400%

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ibkr-trading

[Install]
WantedBy=multi-user.target
